# service/log/run
#!/bin/sh
     exec setuidgid root multilog t s16777215 ./main
